# home/.config/systemd/user/default.target
[Unit]
Description="Default user target, starts a bunch of services"

[Install]
Alias=default.target
