# home/.config/systemd/user/home-assistant.service
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/hass

[Install]
WantedBy=multi-user.target
