# home/.config/systemd/user/rescuetime.service
[Unit]
Description=Rescuetime
Documentation=
Requires=wait-for-xserver.service
After=wait-for-xserver.service

[Service]
Type=simple
ExecStart=/usr/bin/rescuetime
Environment="DISPLAY=:0"
Restart=always

[Install]
WantedBy=default.target
