# home/.config/systemd/user/sync-annex.service
[Unit]
Description=Sync annexes

[Service]
ExecStart=/usr/bin/bash /home/erb/.bin/sync-annex

[Install]
WantedBy=multi-user.target
